# docs/launchd/com.helm.shquant.signals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.helm.shquant.signals</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/helm/Documents/Code/sh_quant/.venv/bin/python</string>
    <string>/Users/helm/Documents/Code/sh_quant/scripts/pull_signals.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/helm/Documents/Code/sh_quant</string>

  <!-- 排在 daily/sector_turnover 之后跑: 个股+板块数据落地后再算信号 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>45</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/helm/.market_data/logs/signals.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/helm/.market_data/logs/signals.err.log</string>
</dict>
</plist>
